# net-analyzer/adeona/files/adeona
#!/sbin/runscript

start() {
  start-stop-daemon --start --background \
      --pidfile /var/run/adeona.pid --make-pidfile \
      --exec /usr/sbin/adeona-client.exe -- \
      -s /var/lib/adeona/adeona-clientstate.cst \
      -r /var/lib/adeona -l /var/log/adeona/
  eend $?
}

stop() {
  start-stop-daemon --stop --exec /usr/sbin/adeona-client.exe \
    --pidfile /var/run/adeona.pid
  eend $?
}
